# crontab
30 17 * * FRI /usr/local/bin/python /news_sentiment/main.py > /proc/1/fd/1 2>&1
